# target/linux/ramips/mt7621/base-files/etc/init.d/bootcount
#!/bin/sh /etc/rc.common

START=99

boot() {
	case $(board_name) in
	alfa-network,quad-e4g)
		[ -n "$(fw_printenv bootcount bootchanged 2>/dev/null)" ] &&\
			echo -e "bootcount\nbootchanged\n" | /usr/sbin/fw_setenv -s -
		;;
	beeline,smartbox-turbo|\
	beeline,smartbox-turbo-plus|\
	rostelecom,rt-sf-1)
		[[ $(hexdump -n 1 -e '/1 "%1d"' -s $((0x20001)) /dev/mtd3) == \
			$((0xff)) ]] || printf '\xff' | dd of=/dev/mtdblock3 \
			count=1 bs=1 seek=$((0x20001))
		;;
	jdcloud,re-cp-02)
		echo -e "bootcount 0\nbootlimit 5\nupgrade_available 1"  | /usr/sbin/fw_setenv -s -
		;;
	linksys,e5600|\
	linksys,ea6350-v4|\
	linksys,ea7300-v1|\
	linksys,ea7300-v2|\
	linksys,ea7500-v2|\
	linksys,ea8100-v1|\
	linksys,ea8100-v2)
		mtd resetbc s_env || true
		;;
	belkin,rt1800|\
	linksys,e7350|\
	samknows,whitebox-v8)
		fw_setenv bootcount 0
		;;
	dna,valokuitu-plus-ex400|\
	genexis,pulse-ex400)
		fw_setenv boot_cnt_primary 0
		fw_setenv boot_cnt_alt 0
		;;
	zyxel,lte3301-plus)
		[ $(printf %d $(fw_printenv -n DebugFlag)) -gt 0 ] || fw_setenv DebugFlag 1
		[ $(printf %d $(fw_printenv -n Image1Stable)) -gt 0 ] || fw_setenv Image1Stable 1
		[ $(printf %d $(fw_printenv -n Image1Try)) -gt 0 ] && fw_setenv Image1Try 0
		;;
	zyxel,lte5398-m904|\
	zyxel,lte7490-m904|\
	zyxel,nr7101)
		[ $(printf %d $(fw_printenv -n DebugFlag)) -gt 0 ] || fw_setenv DebugFlag 0x1
		[ $(printf %d $(fw_printenv -n Image1Stable)) -gt 0 ] || fw_setenv Image1Stable 1
		[ $(printf %d $(fw_printenv -n Image1Try)) -gt 0 ] && fw_setenv Image1Try 0
		;;
	teltonika,rutm51|\
	teltonika,rutm11)
		#this will disable failsafe boot, to recover to factory value execute and flash fw over bootloader
		# fw_setenv bootcmd "bootfailsafe || reset"
		#fw_setenv bootcmd "bootnand 0"
		#find_mtd_part bootconfig-a
		local PREFIX=/dev/mtd
		
		local INDEX="$(find_mtd_index bootconfig-a)"
		local bootconfig_a="$PREFIX$INDEX"
		
		INDEX="$(find_mtd_index bootconfig-b)"
		local bootconfig_b="$PREFIX$INDEX"
		local active_part
		local seek=0
		
		logger "bootconfig-a: $bootconfig_a"
		logger "bootconfig-b: $bootconfig_b"
		
		if grep -q ubi.mtd=rutos-a /proc/cmdline; then
			active_part="$bootconfig_a"
		else
			active_part="$bootconfig_b"
		fi

		while [ $(hexdump -n 16 -s "$seek" -e '16/1 "%02x"' "$active_part") != "ffffffffffffffffffffffffffffffff" ]; do
			seek=$(( seek + 16 ))
		done

		seek=$(( seek - 16 ))

		if grep -q ubi.mtd=rutos-a /proc/cmdline; then
			[ "$(hexdump -n 16 -s "$seek" -e '16/1 "%02x"' "$bootconfig_a")" != "e1b0baba0100f901f801f0003dd4fb95" ] && {
				echo "Updating boot-a config"
				mtd erase "$bootconfig_a"
				echo -ne "\xe1\xb0\xba\xba\x01\x00\xf9\x01\xf8\x01\xf0\x00\x3d\xd4\xfb\x95" | mtd write - "$bootconfig_a"
			}
		else
			[ "$(hexdump -n 16 -s "$seek" -e '16/1 "%02x"' "$bootconfig_b")" != "e1b0baba0101f800f901f000f9420c7d" ] && {
				#~ mtd erase "bootconfig-b"
				echo "Updating boot-b config"
				mtd erase "$bootconfig_b"
				echo -ne "\xe1\xb0\xba\xba\x01\x01\xf8\x00\xf9\x01\xf0\x00\xf9\x42\x0c\x7d" | mtd write - "$bootconfig_b"
			}
		fi
		;;
	esac
}
